# roles/hadoop_ha/templates/namenode.service
[Unit]
Description=Description for sample script goes here
After=journalnode.service

[Service]
User=hadoop
Type=forking
ExecStart=/bin/bash -c "source /home/hadoop/.bash_profile && /usr/bigdata/Envs/scripts/bigadm.py --start namenode --type=service"
ExecStop=/bin/bash -c "source /home/hadoop/.bash_profile && /usr/bigdata/Envs/scripts/bigadm.py --stop namenode --type=service"
TimeoutStartSec=0

[Install]
WantedBy=default.target
